# sources/meta-prusa/recipes-core/udev/udev-usbmount/rmdir-on-exit@.service
[Unit]
PartOf=%i.automount

[Service]
Type=oneshot
ExecStop=-/bin/rmdir --verbose %f
RemainAfterExit=yes
