# sources/meta-prusa/recipes-firmware/sla-slicer-upload/sla-slicer-upload/sla-slicer-upload-restarter.path
[Path]
PathModified=/etc/sl1fw/slicer-upload-api.key

[Install]
WantedBy=multi-user.target
