# recipes-firmware/sl1fw/sl1fw/home-root-ramdisk.mount
[Unit]
Description=Printing tmpfs mount

[Mount]
What=none
Where=/home/root/ramdisk
Type=tmpfs
Options=defaults

[Install]
WantedBy=local-fs.target
