# recipes-firmware/sl1fw/sl1fw/home-root-usb.automount
[Unit]
Description=Automount user flashdisk

[Automount]
Where=/home/root/usb
TimeoutIdleSec=1

[Install]
WantedBy=multi-user.target
