# recipes-firmware/sl1fw/sl1fw/mnt-rootfs.automount
[Unit]
Description=Automount update rootfs

[Automount]
Where=/mnt/rootfs
TimeoutIdleSec=1

[Install]
WantedBy=local-fs.target
