# recipes-firmware/sl1fw/sl1fw/mnt-usb.automount
[Unit]
Description=Automount user flashdisk

[Automount]
Where=/mnt/usb
TimeoutIdleSec=1

[Install]
WantedBy=local-fs.target
