# recipes-core/powerpanic/files/panic.service
[Unit]
Description=Power panic response service

[Service]
ExecStart=/usr/bin/panic
 
[Install]
WantedBy=multi-user.target
