# recipes-firmware/avahi-restarter/files/avahi-restart.path
[Unit]
Description=Start avahi-restart when hostname changes

[Path]
PathChanged=/etc/hostname

[Install]
WantedBy=multi-user.target
