# sources/meta-prusa/recipes-firmware/avahi-restarter/files/avahi-restart.service
[Unit]
Description=Restart avahi

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart avahi-daemon
