# sources/meta-prusa/recipes-bsp/etc/files/boot.mount
[Unit]
Description=Mount boot partition
DefaultDependencies=yes

[Mount]
What=${BOOT_DEV}
Where=/boot
Type=vfat

[Install]
WantedBy=local-fs.target
