# recipes-support/udiskie/udiskie/udiskie.service
[Unit]
Description=Udiskie automount

[Service]
ExecStart=/usr/bin/udiskie -c /etc/udiskie.yaml

[Install]
WantedBy=multi-user.target
